# roles/linux/windowmanager/templates/metacity.service.j2
[Unit]
Description=Metacity service (windowmanager)
After=xvfb.target

[Service]
ExecStart=/usr/bin/metacity --display=:{{screen_num}}

[Install]
WantedBy=multi-user.target
