# roles/linux/xvfb/templates/xvfb.service.j2
[Unit]
Description=X Virtual Frame Buffer Service
After=network.target

[Service]
ExecStart=/usr/bin/Xvfb :{{xvfb_screen_num}} -screen 0 {{xvfb_resolution}} -ac

[Install]
WantedBy=multi-user.target
